# minecraft/files/ddns/duckdns.cron
*/5 * * * * root /var/minecraft/ddns/duckdns.sh >/dev/null 2>&1
@reboot     root /var/minecraft/ddns/duckdns.sh >/dev/null 2>&1
